# subsonic-installer-debian/src/DEBIAN/preinst
#! /bin/sh

set -e

# Stop Supersonic service.
if [ -e /etc/init.d/supersonic ]; then
  invoke-rc.d supersonic stop
fi
